# systemd/displayboard.service.d/30-redisvars.conf
[Service]
Environment=SESSION_TYPE=redis
Environment=SESSION_REDIS=redis://127.0.0.1:6379